# templates/systemd/wireguard-dns.service.j2
[Unit]
Description=Dynamic DNS for WireGuard clients

[Service]
Type=simple
ExecStart={{ outputdir }}/update_dnsmasq_with_wireguard_clients.sh
